# etc/systemd/system/ecowitt-forwarder.service
[Unit]
Description=Ecowitt → WeeWX Interceptor Forwarder
After=network-online.target weewx.service
Wants=network-online.target

[Service]
Type=simple
User=username
WorkingDirectory=/home/username/dev/ecowitt-forwarder
ExecStart=/usr/bin/python3 /home/username/dev/ecowitt-forwarder/ecowitt-forwarder.py
Restart=always
RestartSec=2
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
